# etv.service
[Unit]
Description=etv

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/home/jrwren/go/bin/etv

[Install]
WantedBy=multi-user.target
